# crontabs/return-of-results
SHELL=/bin/bash
PATH=/opt/backoffice/bin:/usr/local/bin:/usr/bin:/bin

# Point Pipenv at the production environment
PIPENV_PIPFILE=/opt/backoffice/id3c-production/Pipfile

LOG_LEVEL=warning
PGSERVICE=production-etl

# Base of our env.d directories
ENVD=/opt/backoffice/id3c-production/env.d


# Generate data file and PDFs for returning SFS results via UW Lab Med's SecureLink portal
5,35 * * * * ubuntu promjob "return of results" flock --no-fork --nonblock /var/lock/return-of-results fatigue --quiet pipenv run envdir $ENVD/redcap-scan/ envdir $ENVD/redcap-sfs envdir $ENVD/securelink/ /opt/backoffice/bin/return-of-results/generate
