# crontabs/scan-return-of-results
SHELL=/bin/bash
PATH=/usr/local/bin:/usr/bin:/bin

# Point Pipenv at the production environment
PIPENV_PIPFILE=/opt/backoffice/id3c-production/Pipfile

LOG_LEVEL=warning
PGSERVICE=production-etl

# Base of our env.d directories
ENVD=/opt/backoffice/id3c-production/env.d


# Generate data file and PDFs for returning SCAN IRB results via UW Lab Med's SecureLink portal
35 */2 * * * ubuntu flock --no-fork --nonblock /var/lock/scan-return-of-results pipenv run envdir $ENVD/redcap-scan/ envdir $ENVD/securelink/ /opt/backoffice/bin/scan-return-of-results/generate
