# crontabs/scan-switchboard
SHELL=/bin/bash
PATH=/opt/backoffice/bin:/usr/local/bin:/usr/bin:/bin

# Point Pipenv at the SCAN Switchboard environment
PIPENV_PIPFILE=/opt/scan-switchboard/Pipfile

# SCAN Switchboard uses the same environment variables as id3c-production
ENVD=/opt/backoffice/id3c-production/env.d

# Configure a private runtime directory for "fatigue".
XDG_RUNTIME_DIR=/home/ubuntu/run


# Refresh the SCAN Switchboard SQLite database every 10 minutes
*/10 * * * * ubuntu chronic fatigue envdir $ENVD/redcap-sfs/ envdir $ENVD/redcap-scan/ pipenv run make -BC /opt/scan-switchboard
