# crontabs/website-science-map
SHELL=/bin/bash
PATH=/usr/local/bin:/usr/bin:/bin

LOG_LEVEL=warning
PGSERVICEFILE=/opt/backoffice/pg_service.conf
PGSERVICE=production-etl

# Export the latest sample collection data to S3 for seattleflu.org/science
28 * * * * ubuntu /opt/backoffice/export/science-map-geojson | aws s3 cp - s3://seattle-flu-study/samples.geojson
